# AppShell/LaunchAgents/com.inceptyon.gargantua.scheduler.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.inceptyon.gargantua.scheduler</string>
    <key>BundleProgram</key>
    <string>Contents/MacOS/GargantuaScheduler</string>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/tmp/gargantua-scheduler.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/gargantua-scheduler.log</string>
</dict>
</plist>
